# .docker/rootfs/etc/supervisord.d/supervisord.ini
[supervisord]
nodaemon=true

[program:xl2tpd]
command=/usr/sbin/xl2tpd -D

[program:strongswan]
command=/usr/sbin/strongswan start --nofork
